# roles/apps/maven/templates/profile.d/maven.sh
#!/bin/bash
export M2_HOME=/opt/maven
export PATH=$PATH:$M2_HOME/bin